# vps/systemd/ai-office-docker.service
[Unit]
Description=AI Office Monitoring Docker Compose Stack
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/ai-office/repo
ExecStartPre=/usr/bin/docker compose config -q
ExecStart=/usr/bin/docker compose up -d grafana prometheus node-exporter loki promtail
ExecStop=/usr/bin/docker compose stop grafana prometheus node-exporter loki promtail
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
