# vps/systemd/ai-office-new-ui.service
[Unit]
Description=AI Office — OpenClaw Office UI (BigTech изометрический UI)
After=network.target openclaw.service

[Service]
Type=simple
ExecStart=/usr/bin/openclaw-office --port 3001
Restart=on-failure
RestartSec=5
Environment=OPENCLAW_GATEWAY_URL=wss://80.74.25.43/gateway-ws
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
