# vps/systemd/ai-office-postgres-firewall.service
[Unit]
Description=Block public access to AI Office Postgres Docker port
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ai-office-postgres-firewall.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
