# vps/systemd/ai-office-docs.service
[Unit]
Description=AI Office Documentation (Docusaurus)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/serve /opt/ai-office/docs/build -p 5000 -n
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
