# vps/systemd/openclaw.service
[Unit]
Description=OpenClaw AI Office Gateway
After=network-online.target ollama.service
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/root
EnvironmentFile=/root/.env
Environment=OLLAMA_HOST=http://localhost:11434
Environment=POSTGRES_HOST=localhost
Environment=POSTGRES_PORT=5432
Environment=POSTGRES_DB=ai_office
Environment=POSTGRES_USER=postgres
ExecStartPre=/usr/local/sbin/openclaw-preflight.sh
ExecStart=/usr/bin/openclaw gateway run --port 18789 --bind lan --ws-log compact
Restart=always
RestartSec=5s
TimeoutStartSec=90s
TimeoutStopSec=30s
KillMode=control-group
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
