# vps/systemd/openclaw-watchdog.service
[Unit]
Description=OpenClaw watchdog health check
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/openclaw-watchdog.sh
